# env.service
[Unit]
Description=environment for ToyVpnServer
After=network.target

[Service]
ExecStart=/home/centos/src/ToyVpnServer/env.sh

[Install]
WantedBy=multi-user.target
